# crontab
* * * * * ubuntu /bin/bash /home/ubuntu/apps/barkedu/repository/run_on_server.sh fab $DEPLOYMENT_TARGET cron_jobs.test >/var/log/barkedu/crontab.log 2>&1
